# DNI/package/lan-scan/files/lan-scan.init
#!/bin/sh /etc/rc.common

start() {
	/usr/sbin/lan-scan
}

stop() {
	killall lan-scan
}

boot() {
	mkdir -p /tmp/lanscan
	start
}
